# var/lib/dpkg/info/nftables.list
/.
/etc
/etc/nftables.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/nftables.service
/usr
/usr/sbin
/usr/sbin/nft
/usr/share
/usr/share/doc
/usr/share/doc/nftables
/usr/share/doc/nftables/README.Debian
/usr/share/doc/nftables/changelog.Debian.gz
/usr/share/doc/nftables/copyright
/usr/share/doc/nftables/examples
/usr/share/doc/nftables/examples/README
/usr/share/doc/nftables/examples/all-in-one.nft
/usr/share/doc/nftables/examples/arp-filter.nft
/usr/share/doc/nftables/examples/bridge-filter.nft
/usr/share/doc/nftables/examples/ct_helpers.nft
/usr/share/doc/nftables/examples/inet-filter.nft
/usr/share/doc/nftables/examples/inet-nat.nft
/usr/share/doc/nftables/examples/ipv4-filter.nft
/usr/share/doc/nftables/examples/ipv4-mangle.nft
/usr/share/doc/nftables/examples/ipv4-nat.nft
/usr/share/doc/nftables/examples/ipv4-raw.nft
/usr/share/doc/nftables/examples/ipv6-filter.nft
/usr/share/doc/nftables/examples/ipv6-mangle.nft
/usr/share/doc/nftables/examples/ipv6-nat.nft
/usr/share/doc/nftables/examples/ipv6-raw.nft
/usr/share/doc/nftables/examples/load_balancing.nft
/usr/share/doc/nftables/examples/nat.nft
/usr/share/doc/nftables/examples/netdev-ingress.nft
/usr/share/doc/nftables/examples/overview.nft
/usr/share/doc/nftables/examples/pf.os
/usr/share/doc/nftables/examples/secmark.nft
/usr/share/doc/nftables/examples/sets_and_maps.nft
/usr/share/doc/nftables/examples/sysvinit
/usr/share/doc/nftables/examples/sysvinit/README
/usr/share/doc/nftables/examples/sysvinit/nftables.init
/usr/share/doc/nftables/examples/workstation.nft
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/nft.8.gz
/usr/share/man/man8/nftables.8.gz
